# openwrt/attitude_adjustment/package/base-files/files/etc/crontabs/root
*/30 * * * * /usr/local/bin/beacon
0,30 * * * * cat /var/run/latlon.js >/tmp/latlon.js
*/5 * * * * /usr/local/bin/olsrd_watchdog
*/5 * * * * killall -HUP dnsmasq
